# splunk_rules/scheduled_task_created.yml
title: Scheduled Task Created
id: 40e1f2a3-b4c5-6789-defa-012345678901
status: stable
description: Yeni scheduled task yaradılmasını aşkarlayır (T1053.005 - persistence).
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: SOC team
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4698
  condition: selection
level: medium
tags:
  - attack.persistence
  - attack.t1053.005
